# spinmaster_main.service
[Unit]
Description=Test passing multiple arguments

[Service]
#Environment="SCRIPT_ARGS=%I"
ExecStart=echo %i