# spinmaster_main@.service
[Unit]
Description=Test passing multiple arguments

[Service]
Environment="SCRIPT_ARGS=%i"
ExecStart=echo ${SCRIPT_ARGS}