# run/spinmaster_main@.service
[Unit]
Description=Test passing multiple arguments

[Service]
EnvironmentFile=./spinmaster_main_args
ExecStart=echo